# htcondor-ce/deploy/docker/config/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:condor]
command=/usr/sbin/condor_master
autostart=true
autorestart=true
